# k2fsa/rootfs/etc/s6-overlay/s6-rc.d/k2fsa/run
#!/command/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# Start Whisper service
# ==============================================================================





if [ ! -f "/module/start.sh" ]; then
    cp /k2fsa/start.sh "/module/start.sh"
fi

/module/start.sh