# k2fsa/rootfs/etc/s6-overlay/s6-rc.d/k2fsa/run
#!/command/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# Start Whisper service
# ==============================================================================

if [ "$(uname -m)" == "x86_64" ] && ! grep -qw 'avx' /proc/cpuinfo; then
    bashio::log.warning "Your CPU does not support Advanced Vector Extensions (AVX). Whisper will run slower than normal."
fi

model="$(bashio::config 'model')"
if [ "${model}" = 'custom' ]; then
    # Override with custom model
    model="$(bashio::config 'custom_model')"
    if [ -z "${model}" ]; then
      bashio::exit.nok "Custom model is not set"
    fi
fi

if bashio::config.true 'debug_logging'; then
    flags+=('--debug')
fi



if [ ! -f "/module/start.sh" ]; then
    cp /k2fsa/start.sh "/module/start.sh"
fi

/module/start.sh