# k2fsa-python/rootfs/etc/s6-overlay/s6-rc.d/k2fsa/run
#!/command/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
# ==============================================================================
# Start Whisper service
# ==============================================================================

while true; do
    echo "hello world"
    sleep 5
done