# OpenWrt/overlay_setup
cat << 'EOF' > /etc/init.d/overlay_setup
#!/bin/sh /etc/rc.common
# OverlayFS Setup Script
START=00

start() {
    echo "Setting up OverlayFS for /"

    mkdir -p /data/upper

    # Mount OverlayFS for /
    mount -t overlay overlay -o lowerdir=/,upperdir=/data/upper,workdir=/data/upper /


EOF
