# OpenWrt/overlay_setup
cat << 'EOF' > /etc/init.d/overlay_root
#!/bin/sh /etc/rc.common
# OverlayFS Root Setup Script
# 初始化优先级
START=00

start() {
    echo "Initializing OverlayFS for root filesystem..."

    # 定义目录路径
    DATA_DIR="/data"
    UPPER_DIR="$DATA_DIR/upper"
    WORK_DIR="$DATA_DIR/work"

    # 确保持久化存储目录存在
    mkdir -p $UPPER_DIR $WORK_DIR

    # 检查是否已经迁移过
    if [ ! -f "$UPPER_DIR/.migration_complete" ]; then
        echo "Migrating root filesystem to $UPPER_DIR..."
        
        # 使用 rsync 迁移文件系统，排除临时目录和持久化目录
        rsync -a --exclude="$DATA_DIR" --exclude="/proc" --exclude="/sys" --exclude="/dev" / $UPPER_DIR/
        
        # 标记迁移完成
        touch "$UPPER_DIR/.migration_complete"
        echo "Migration completed."
    else
        echo "Migration already completed. Skipping."
    fi

    # 挂载 OverlayFS
    echo "Setting up OverlayFS..."
    mount -t overlay overlay -o lowerdir=/,upperdir=$UPPER_DIR,workdir=$WORK_DIR /

    # 验证挂载成功
    if mount | grep -q "overlay on / type overlay"; then
        echo "OverlayFS setup successful."
    else
        echo "OverlayFS setup failed!" >&2
        exit 1
    fi
}

EOF
