# Alist/rootfs/etc/services.d/alist/run
#!/usr/bin/with-contenv bashio
set -e

if bashio::config.true 'ssl'; then
    bashio::log.info "SSL is enabled. Copying SSL certs..."
    mkdir -p /ssl
    cp "$(bashio::config 'certfile')" /ssl/fullchain.pem
    cp "$(bashio::config 'keyfile')" /ssl/privkey.pem
fi

bashio::log.info "Starting AList..."
exec /alist server